# kaabitseja/crontab
*/5 * * * * /usr/local/bin/python3 /app/main.py >> /var/log/cron.log 2>&1
