# kaabitseja/crontab
*/5 * * * * /usr/bin/java -cp /app Main "/app/data/" >> /var/log/cron.log 2>&1
